# files/etc/profile.d/01-golang.sh
export GOROOT=/usr/lib/go
export GOPATH=~/go
export GO111MODULE=auto