# files/etc/supervisor.d/code-server.ini
[program:code-server]
command=sudo -u www-data /usr/lib/code-server/bin/code-server --bind-addr 0.0.0.0:8000 --user-data-dir /data/home/code --disable-telemetry --disable-update-check
priority=100
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
